# php-worker/supervisord.d/base_be_queue_worker.conf
[program:base-be-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/base_be/artisan queue:work -v --max-time=3600 --timeout=60 --sleep=3 --tries=3 --backoff=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/base_be/storage/logs/queue/process%(process_num)s.log
stderr_logfile=/var/www/base_be/storage/logs/queue/process%(process_num)s.error.log
stopwaitsecs=3600
